# services/reaction-listener/reaction-listener.service
[Unit]
Description=autoevolve Discord Reaction Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=SYSTEM_USER
Group=SYSTEM_USER

# Set these to match your agent's setup
Environment=BOT_TOKEN_PATH=/path/to/.openclaw/openclaw.json
Environment=SIGNALS_PATH=/opt/autoevolve/local/signals.jsonl

ExecStart=/usr/bin/python3 /opt/autoevolve/services/reaction-listener/listener.py
WorkingDirectory=/opt/autoevolve/services/reaction-listener

Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autoevolve-reactions

[Install]
WantedBy=multi-user.target
